# mirrorlist-timer/update_mirrorlist.timer
[Unit]
Description=Automatic update of the pacman mirrorlist

[Timer]
OnCalendar=weekly UTC
Persistent=true

[Install]
WantedBy=timers.target
