# mirrorlist-timer/update_mirrorlist.service
[Unit]
Description=Automatic update of the pacman mirrorlist

[Service]
Type=oneshot
ExecStart=generate_mirrorlist.sh
